# systemd/bluetooth-auto-restore.service
[Unit]
Description=Restore Bluetooth pairings from backup after reboot
After=bluetooth.service
Wants=bluetooth.service

[Service]
Type=oneshot
ExecStart=/usr/local/bin/bluetooth-auto-restore
RemainAfterExit=yes

[Install]
WantedBy=multi-user.target
