# system/crontab/crontab.pi
# Pi user crontab (reference only — install with: crontab crontab.pi)
# Minetest world backup: weekly Sunday 3am
0 3 * * 0 /usr/bin/python3 /home/whf/dotfiles/minetest-backup/backup_minetest.py --config /home/whf/dotfiles/minetest-backup/config.toml >/dev/null 2>&1
